# conky.service
[Unit]
Description=Run conky

[Service]
Type=simple
ExecStart=/usr/bin/conky

[Install]
WantedBy=default.target
